# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:fpm]
command=/usr/sbin/php5-fpm -F